# core/system-config/autostart/cosmic-ai.desktop
[Desktop Entry]
Type=Application
Name=Cosmic AI
GenericName=AI Assistant
Comment=Cosmic OS AI-powered assistant
Exec=/usr/bin/python3 -m core.ai_engine.main
Icon=cosmic-ai
Terminal=false
Categories=System;Utility;
StartupNotify=false
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
X-KDE-autostart-phase=2
Hidden=false
NoDisplay=false
