# core/system-config/autostart/cosmic-ai.desktop
[Desktop Entry]
Type=Application
Name=Cosmic AI
Comment=AI Assistant for Cosmic OS
Exec=/usr/bin/python3 /opt/cosmic-os/core/ai_engine/main.py
Icon=applications-science
Terminal=false
Categories=System;Utility;
StartupNotify=false
X-KDE-autostart-phase=2
X-GNOME-Autostart-enabled=true
OnlyShowIn=KDE;GNOME;XFCE;
