# core/system-config/systemd/cosmic-ai.service
[Unit]
Description=Cosmic OS AI Assistant Daemon
Documentation=https://github.com/cosmic-os/cosmic-os
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/cosmic-os/core/ai_engine/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Environment
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cosmic-ai

# Security (optional hardening)
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=%h/.config/cosmic-os %h/.local/share/cosmic-os /tmp

[Install]
WantedBy=graphical-session.target
