# core/system-config/systemd/cosmic-hotreload.service
[Unit]
Description=Cosmic OS Hot Reload Service
Documentation=https://github.com/cosmic-os/cosmic-os
After=cosmic-ai.service
Wants=cosmic-ai.service

[Service]
Type=simple
ExecStart=/opt/cosmic-os/scripts/dev-sync.sh
Restart=on-failure
RestartSec=3

# Environment
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Watch the shared folder (configure path as needed)
Environment=COSMIC_WATCH_DIR=/mnt/cosmic-os

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cosmic-hotreload

[Install]
WantedBy=graphical-session.target
